# pi-shutdown.service
[Unit]
Description=Tiny Flask API to shut down the Pi
After=network.target

[Service]
User=root
# optional token override
EnvironmentFile=/etc/default/pi-shutdown
WorkingDirectory=/opt/pi-shutdown
ExecStart=/usr/bin/python3 -m gunicorn -w 2 -b 0.0.0.0:8000 api:app
# Let the shutdown command run even if Flask is killed
KillMode=process
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
